# service/log/run
#!/bin/sh
# Runit log service execution file using multilog.

LOG_DIR="/var/log/gps_socat"

# Merge stderr into stdout (2>&1)
exec 2>&1
# Execute multilog command
exec multilog t s25000 n4 $LOG_DIR
